# docker/etc/s6/bitbucket/run
#!/bin/bash

# 启动X Server服务器
exec gosu "${USERNAME}" /opt/atlassian/bitbucket/bin/start-bitbucket.sh -fg --no-search
