# docker/etc/s6/bitbucket/run
#!/bin/bash

# 启动Bitbucket服务器
exec gosu "${USERNAME}" /opt/atlassian/bitbucket/bin/start-bitbucket.sh --no-search
